# emailbot.timer
[Unit]
Description=Danapfel Digital E-Mail Bot Timer (Alle 2 Wochen Montags 08:00)
Requires=emailbot.service

[Timer]
# Wird jeden Montag um 08:00 Uhr ausgelöst.
# Das Skript selbst sorgt dafür, dass nur alle 14 Tage wirklich gesendet wird.
OnCalendar=Mon *-*-* 08:00:00

# Persistenz: Wenn System zur geplanten Zeit aus war, beim nächsten Start nachholen
Persistent=true

# Genauigkeit (optional, Standard ist 1 Minute)
AccuracySec=1min

[Install]
WantedBy=timers.target

