# emailbot.service
[Unit]
Description=Danapfel Digital EmailBot
After=network.target

[Service]
Type=simple
User=exit
WorkingDirectory=/home/exit/emailbot
ExecStart=/usr/bin/node /home/exit/emailbot/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=emailbot

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
